# deploy/termigate-container.service
# termigate — Container deployment (Podman)
#
# Install:
#   sudo useradd -r -m -d /var/lib/termigate -s /usr/sbin/nologin termigate
#   sudo mkdir -p /etc/termigate /var/lib/termigate
#   sudo chown -R termigate:termigate /var/lib/termigate
#   # Build or pull the image:
#   sudo podman build -t localhost/termigate:latest -f Containerfile .
#   sudo cp deploy/termigate-container.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now termigate-container.service
#
# By default this runs tmux *inside* the container. To connect to the host's
# tmux sessions instead, add a bind mount for the tmux socket directory and
# set TERMIGATE_TMUX_SOCKET (see Containerfile docs in README.md).

[Unit]
Description=termigate — Browser-based tmux session manager (container)
Documentation=https://github.com/anthropics/termigate
After=network-online.target
Wants=network-online.target
Requires=podman.service

[Service]
Type=exec
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

Environment=TERMIGATE_IMAGE=localhost/termigate:latest
Environment=TERMIGATE_PORT=8888
# TERMIGATE_SECRET_KEY_BASE is required and intentionally not shipped here.
# Generate one with `openssl rand -base64 48` and supply it via a drop-in:
#   sudo systemctl edit termigate-container.service
# adding:
#   [Service]
#   Environment=TERMIGATE_SECRET_KEY_BASE=<generated-value>
# The server will refuse to start without a real value.
#
# If termigate sits behind a reverse proxy (Caddy, nginx, Traefik, etc.),
# set TERMIGATE_TRUSTED_PROXIES to a comma-separated list of proxy CIDRs
# so X-Forwarded-For is honored. Without this, rate limits and auth logs
# collapse all traffic onto the proxy IP. Example:
#   Environment=TERMIGATE_TRUSTED_PROXIES=127.0.0.0/8,::1/128

# Remove any stale container before starting
ExecStartPre=-/usr/bin/podman rm -f termigate

ExecStart=/usr/bin/podman run \
    --rm \
    --name termigate \
    --publish ${TERMIGATE_PORT}:8888 \
    --env PHX_SERVER=true \
    --env PHX_HOST=localhost \
    --env SECRET_KEY_BASE=${TERMIGATE_SECRET_KEY_BASE} \
    --env TERMIGATE_TRUSTED_PROXIES=${TERMIGATE_TRUSTED_PROXIES:-} \
    --volume /etc/termigate:/etc/termigate:ro \
    --volume termigate-data:/var/lib/termigate \
    ${TERMIGATE_IMAGE}

ExecStop=/usr/bin/podman stop --time 10 termigate

[Install]
WantedBy=multi-user.target
