# packaging/scriptlets/postinst
#!/bin/sh
# postinst — set ownership, reload systemd, notify admin.

set -e

chown -R termigate:termigate /var/lib/termigate /usr/lib/termigate
chmod 750 /var/lib/termigate

if command -v systemctl > /dev/null 2>&1; then
    systemctl daemon-reload || :
fi

cat <<'EOF'

termigate is installed but not enabled.

Next steps:
  1. Generate a secret key:
       openssl rand -base64 48
  2. Add it via a systemd drop-in (the shipped service file does NOT include
     a placeholder — the server refuses to start without a real value):
       sudo systemctl edit termigate.service
     and paste:
       [Service]
       Environment=SECRET_KEY_BASE=<generated-value>
  3. Optionally copy the sample config:
       cp /etc/termigate/config.example.yaml /etc/termigate/config.yaml
  4. Enable and start:
       systemctl daemon-reload
       systemctl enable --now termigate

EOF

exit 0
